# staging/archive/old_scripts/default.service
[Unit]
Description=desc_here

[Service]
ExecStart=path_name_here

[Install]
WantedBy=multi-user.target